# dist/supervisor.conf
[program:coldstar]
command=twistd --pidfile coldstar.pid -n coldstar --config /srv/coldstar/code/local.conf
stdout_logfile=/srv/coldstar/logs/coldstar.log
stderr_logfile=/srv/coldstar/logs/coldstar.log
#user=coldstar
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
directory=/srv/coldstar/code
environment=PYTHONUNBUFFERED=1, PYTHONPATH=/srv/coldstar/code/src, VIRTUAL_ENV=/srv/coldstar/.virtualenv, PATH="/srv/coldstar/.virtualenv/bin:$PATH", PYTHONHOME=/srv/coldstar/.virtualenv


[program:simargl]
command=twistd --pidfile simargl.pid -n coldstar --config /srv/coldstar/code/local_simargl.conf
stdout_logfile=/srv/coldstar/logs/simargl.log
stderr_logfile=/srv/coldstar/logs/simargl.log
#user=coldstar
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
directory=/srv/coldstar/code
environment=PYTHONUNBUFFERED=1, PYTHONPATH=/srv/coldstar/code/src, VIRTUAL_ENV=/srv/coldstar/.virtualenv, PATH="/srv/coldstar/.virtualenv/bin:$PATH", PYTHONHOME=/srv/coldstar/.virtualenv
